# packages/agent/scripts/tokentimer-agent.service
# TokenTimer Agent systemd unit (installed by install-agent.sh).
#
# Layout assumed by this unit:
#   /opt/tokentimer-agent/app    agent package (read-only at runtime)
#   /opt/tokentimer-agent/state  config dir: config.json, credential,
#                                signing-key-pin.json, replay-store.json, keys/
#
# The agent is outbound-only (no listening sockets), runs as the unprivileged
# tokentimer-agent user, and only ever writes inside its state dir, so the
# hardening below can be strict without breaking anything.

[Unit]
Description=TokenTimer CertOps Agent (outbound-only execution plane)
Documentation=https://github.com/tokentimer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tokentimer-agent
Group=tokentimer-agent
Environment=TOKENTIMER_AGENT_CONFIG_DIR=/opt/tokentimer-agent/state
# Optional bootstrap env file written by install-agent.sh for first-run
# registration (contains TOKENTIMER_AGENT_BOOTSTRAP_TOKEN, mode 0600).
# The token is single-use; the agent never persists it, and deletes this
# file itself right after a successful registration.
EnvironmentFile=-/opt/tokentimer-agent/state/bootstrap.env
ExecStart=/usr/bin/env node /opt/tokentimer-agent/app/bin/tokentimer-agent.js
Restart=always
RestartSec=5
# The agent exits with this status when the control plane retires it
# (heartbeat "retired" signal). A retired agent must NOT be respawned into
# a permanent 410 loop; keep in sync with AGENT_RETIRED_EXIT_CODE in
# packages/agent/src/index.js.
RestartPreventExitStatus=86
SuccessExitStatus=86

# Hardening. ProtectSystem=strict makes the whole filesystem read-only for
# the service except the explicit ReadWritePaths (its state dir).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tokentimer-agent/state
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
# Outbound HTTPS only; the agent accepts no inbound connections.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
